# files/etc/init.d/enhance
#!/bin/sh /etc/rc.common
# Copyright (C) 2008 OpenWrt.org

START=99

start() {
	#Small Packets
	echo "2 30" >/sys/kernel/debug/hnat/hnat_setting
	echo "4 30" >/sys/kernel/debug/hnat/hnat_setting
}
